# dot_config/systemd/user/randombg.timer
[Unit]
Description=Random wallpaper changer for GNOME

[Timer]
OnUnitActiveSec=1h
Unit=randombg.service

[Install]
WantedBy=timers.target
